# files/system/usr/lib/systemd/system/ninjarmm-agent-build.path
[Unit]
Description=Watch for NinjaRMM Agent RPM and trigger build

[Path]
PathExists=/etc/containers/systemd/ninjarmm/ninjarmm-agent.rpm
Unit=ninjarmm-agent-build.service

[Install]
WantedBy=multi-user.target
